# systemd/discgolf.service
[Unit]
Description=Discord Discgolf Bot
After=network.target

[Service]
User=riks
Environment=HOME=/home/riks
WorkingDirectory=/home/riks/Riksdagen-Backend/
ExecStart=/bin/bash -c "source /home/riks/.nvm/nvm.sh && systemd/discgolf.sh"
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
